# etc/udev/rules.d/52-rndis.rules
ACTION=="add", SUBSYSTEM=="net", DRIVERS=="rndis_host", RUN+="/sbin/dhcpcd $env{INTERFACE}"
