# etc/udev/rules.d/90-arduino.rules
SUBSYSTEM=="tty", ATTRS{idVendor}=="2341", ATTRS{idProduct}=="0043",  RUN+="/bin/bash -c 'stty -F %N -icrnl'"
